# infrastructure/roles/teamcity/templates/teamcity.service.j2
[Unit]
Description=JetBrains TeamCity
Requires=network.target
After=syslog.target network.target

[Service]
Type=forking
EnvironmentFile=/etc/sysconfig/teamcity
ExecStart={{ home_dir }}/default/bin/teamcity-server.sh start
ExecStop={{ home_dir }}/default/bin/teamcity-server.sh stop
User={{ user_name }}

[Install]
WantedBy=multi-user.target
