# infrastructure/roles/server/templates/website.service.j2
[Unit]
Description=Website

[Service]
ExecStart={{ bin_dir }}/start.sh
User={{ user_name }}
WorkingDirectory={{ work_dir }}
TimeoutSec=60

[Install]
WantedBy=multi-user.target